# .config/autostart/polkit-kde-authentication-agent-1.desktop
[Desktop Entry]
Exec=/usr/lib/libexec/polkit-kde-authentication-agent-1
Name=polkit-kde-authentication-agent-1
Type=Application
