# .config/autostart/polkitd.desktop
[Desktop Entry]
Exec=/usr/lib/polkit-1/polkitd
Name=polkitd
Type=Application
